# scripts/voiceflow.service
[Unit]
Description=VoiceFlow - Local Voice Input Assistant
After=graphical-session.target

[Service]
Type=simple
WorkingDirectory=/home/sikm/Project/AI-Project/TTS
ExecStart=/bin/bash -c 'DISPLAY=$(ls /tmp/.X11-unix/ | head -1 | sed s/X/:/) XAUTHORITY=/run/user/$(id -u)/gdm/Xauthority exec /home/sikm/Project/AI-Project/TTS/.venv/bin/python main.py'
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
